# ansible/roles/php/files/configuration/etc/profile.d/phpfarm.sh
export PATH=$PATH:/opt/phpfarm/inst/bin:/opt/phpfarm/inst/current/bin

